# oem/gaming-oem-kde6/desktop-overlay/usr/lib/systemd/system/manjaro-post-install.service
[Unit]
Description=Manjaro Post Install.

[Service]
Type=simple
ExecStart=/usr/bin/manjaro-post-install

[Install]
WantedBy=default.target
